# systemd-units/xinit-login.service
[Unit]
Description=Xinit

[Service]
ExecStart=/bin/su mrk3004 -l -c xinit
Restart=no

[Install]
WantedBy=graphical.target
