# systemd-mounts/boot.automount
[Unit]
Description=Boot

[Mount]
Where=/boot

[Install]
WantedBy=local-fs.target
