# etc/systemd/system/dropcaches.timer
[Unit]
Description=Drop caches periodically

[Timer]
OnBootSec=4h
OnUnitActiveSec=4h

[Install]
WantedBy=basic.target
