# etc/systemd/system/leds.service
[Unit]
Description=Keyboard leds configuration

[Service]
Type=oneshot
ExecStart=/usr/bin/tleds -d 20 -c -q web
ExecStart=/usr/bin/cpuled console -c nfn -f
KillMode=process
Restart=no

[Install]
WantedBy=multi-user.target
