# etc/systemd/system/noip.service
[Unit]
Description=No-IP

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/noip2

[Install]
WantedBy=multi-user.target
