# etc/systemd/system/startx.service
[Unit]
After=systemd-user-sessions.service

[Service]
User=lara
ExecStart=/usr/bin/startx -- vt8
RestartSec=3s

[Install]
WantedBy=multi-user.target
